# backend/scripts/systemd/dataflywheel-judge.timer
[Unit]
Description=Run Farm Manager DataFlywheel Judge daily at 02:00

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
Unit=dataflywheel-judge.service

[Install]
WantedBy=timers.target
